# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/parking/supervisord.log
logfile_maxbytes=10MB
loglevel=info

[program:mysql]
command=/usr/bin/mysqld_safe
autostart=true
autorestart=true
stdout_logfile=/var/log/parking/mysql.log
stderr_logfile=/var/log/parking/mysql.log
priority=10

[program:algorithm]
command=python3 /app/algorithm/app.py
directory=/app/algorithm
autostart=true
autorestart=true
stdout_logfile=/var/log/parking/algorithm.log
stderr_logfile=/var/log/parking/algorithm.log
priority=20
environment=KMP_DUPLICATE_LIB_OK="TRUE"

[program:backend]
command=java -jar /app/backend/app.jar
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/var/log/parking/backend.log
stderr_logfile=/var/log/parking/backend.log
priority=30

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/parking/nginx.log
stderr_logfile=/var/log/parking/nginx.log
priority=40
